# fixes/indicator-session/craftparking-force-restart-menuitem.conf
description "Craftparking: force the Restart menu item to show in the system dropdown"

# indicator-session hides Restart whenever it sees a D-Bus name owner for the old
# com.canonical.Unity.Session.EndSessionDialog interface, assuming Unity will show its own
# combined shutdown/reboot prompt - which our Dialogs.qml fix no longer does, so Restart
# disappears. force-restart-menuitem is indicator-session's own built-in override for exactly
# this situation.

start on indicator-services-start
task

exec gsettings set com.canonical.indicator.session force-restart-menuitem true
