# fixes/wireless-display/enable-widi.conf
description "Enable wireless display (widi) support flag once the Android container's property service is confirmed up"

# must run after "android", not inside mount-android.conf - the property
# service isn't up yet at that point and setprop hangs the whole boot
start on android

task

exec setprop ubuntu.widi.supported 1
