# fixes/apt-holds/craftparking-apt-holds.conf
description "hold packages that break boot if upgraded/reconfigured (see apt-dpkg-fix notes)"

# apt-mark hold on an already-held package is a no-op, safe to run every boot.
# root is ro at this point in boot, dpkg needs to write its status file, so
# flip it rw just for this and back to ro right after - same as any other
# fix in this build that touches /etc or dpkg state.
start on filesystem

task

script
    mount -o remount,rw /
    apt-mark hold account-polld address-book-app aethercast aethercast-tools apparmor-easyprof-ubuntu base-files biometryd-bin bluebinder bluez bluez-obexd buteo-syncfw click click-apparmor dialer-app distro-info distro-info-data gir1.2-click-0.4 gstreamer1.0-hybris gstreamer1.0-plugins-bad-faad gstreamer1.0-plugins-bad-videoparsers hfd-service hfd-service-tools history-service isc-dhcp-client isc-dhcp-common isc-dhcp-server libbiometryd1 libbluetooth3 libbuteosyncfw5-0 libclick-0.4-0 libdbus-cpp5 libdeviceinfo0 libgbinder libgbinder-radio libgeonames-common libgeonames0 libgrilio-binder libgsettings-qt1 libgstreamer-plugins-bad1.0-0 libhistoryservice0 libncicore libnciplugin libnet-cpp2 libofono-qt1 libqt5concurrent5 libqt5core5a libqt5dbus5 libqt5feedback5-hfd libqt5gui5 libqt5network5 libqt5opengl5 libqt5printsupport5 libqt5sql5 libqt5sql5-sqlite libqt5test5 libqt5webengine-data libqt5webengine5 libqt5webenginecore5 libqt5widgets5 libqt5xml5 libsystemsettings1 libubuntu-location-service3 libubuntugestures5 libubuntumetrics5 libubuntutoolkit5 libwds0 lxc-android-config messaging-app mir-client-platform-android-caf5 mir-client-platform-android5 mir-graphics-drivers-android mir-graphics-drivers-android-caf mir-platform-graphics-android-caf15 mir-platform-graphics-android15 mobile-broadband-provider-info morph-browser morph-webapp-container nfcd nfcd-binder-plugin nuntium ofono ofono-ril-binder-plugin ofono-scripts packagekit-plugin-click pulseaudio-modules-droid-24 pulseaudio-modules-droid-28 pulseaudio-modules-droid-30 python3-apparmor-click python3-click-package python3-distro-info python3-software-properties qml-module-biometryd qml-module-morph-web qml-module-qtwebengine qml-module-ubuntu-components qml-module-ubuntu-components-labs qml-module-ubuntu-layouts qml-module-ubuntu-performancemetrics qml-module-ubuntu-test qml-module-ubuntu-web-compat qtdeclarative5-gsettings1.0 qtdeclarative5-qtmir-plugin qtdeclarative5-ubuntu-addressbook0.1 qtdeclarative5-ubuntu-contacts0.1 qtdeclarative5-ubuntu-history0.1 qtdeclarative5-ubuntu-keyboard-extensions0.1 qtdeclarative5-ubuntu-telephony-phonenumber0.1 qtdeclarative5-ubuntu-telephony0.1 qtdeclarative5-ubuntu-ui-toolkit-plugin qtmir-android software-properties-common system-image-cli system-image-common system-image-dbus telephony-service ubports-qa-scripts ubuntu-advantage-tools ubuntu-keyboard ubuntu-keyboard-arabic ubuntu-keyboard-azerbaijani ubuntu-keyboard-belarusian ubuntu-keyboard-bosnian ubuntu-keyboard-bulgarian ubuntu-keyboard-catalan ubuntu-keyboard-chinese-chewing ubuntu-keyboard-chinese-pinyin ubuntu-keyboard-croatian ubuntu-keyboard-czech ubuntu-keyboard-danish ubuntu-keyboard-data ubuntu-keyboard-dutch ubuntu-keyboard-emoji ubuntu-keyboard-english ubuntu-keyboard-english-dvorak ubuntu-keyboard-finnish ubuntu-keyboard-french ubuntu-keyboard-german ubuntu-keyboard-greek ubuntu-keyboard-hebrew ubuntu-keyboard-hungarian ubuntu-keyboard-icelandic ubuntu-keyboard-italian ubuntu-keyboard-japanese ubuntu-keyboard-korean ubuntu-keyboard-latvian ubuntu-keyboard-lithuanian ubuntu-keyboard-macedonian ubuntu-keyboard-norwegian-bokmal ubuntu-keyboard-persian ubuntu-keyboard-polish ubuntu-keyboard-portuguese ubuntu-keyboard-romanian ubuntu-keyboard-russian ubuntu-keyboard-scottish-gaelic ubuntu-keyboard-serbian ubuntu-keyboard-slovenian ubuntu-keyboard-spanish ubuntu-keyboard-swedish ubuntu-keyboard-swiss-french ubuntu-keyboard-thai ubuntu-keyboard-turkish ubuntu-keyboard-ukrainian ubuntu-location-service-bin ubuntu-sdk-libs ubuntu-system-settings ubuntu-touch ubuntu-touch-core ubuntu-touch-session ubuntu-touch-session-mir ubuntu-ui-toolkit-theme unity8 unity8-common unity8-private unity8-schemas usensord >/var/log/craftparking-apt-holds.log 2>&1 || true
    mount -o remount,ro /
end script
